# deploy/caravan-post.service
[Unit]
Description=Caravan Post Telegram Bot
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/caravan-post_bot
EnvironmentFile=/var/www/caravan-post_bot/.env
ExecStart=/var/www/caravan-post_bot/.venv/bin/python manage.py runbot
Restart=always
RestartSec=5
TimeoutStopSec=30

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
